# rc.local
#!/bin/bash
#####################################
source global.sh

log $script_name

#####################################